# etc/network/if-pre-up.d/iptables_restore
#!/bin/sh
/sbin/iptables-restore < /etc/iptables.rule
exit