# etc/network/if-up.d/routes
#!/bin/bash
/sbin/route add -net 192.168.230.0 netmask 255.255.255.0
exit